# deploy/lxc/systemd/mikrotik-nms-backend.service
[Unit]
Description=MikroTik NMS backend
Documentation=https://github.com/c0de-ch/mikrotik-c0de
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=mikrotik-nms
Group=mikrotik-nms
WorkingDirectory=/var/lib/mikrotik-nms
EnvironmentFile=/etc/mikrotik-nms/env
ExecStart=/opt/mikrotik-nms/bin/mikrotik-nms
Restart=on-failure
RestartSec=5
LimitNOFILE=65536

# Hardening
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
ProtectClock=true
ProtectHostname=true
ProtectProc=invisible
RestrictSUIDSGID=true
RestrictNamespaces=true
RestrictRealtime=true
LockPersonality=true
MemoryDenyWriteExecute=true
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM
ReadWritePaths=/var/lib/mikrotik-nms

[Install]
WantedBy=multi-user.target
