# deploy/webhook-agent/systemd/mikrotik-nms-deploy-agent.service
[Unit]
Description=MikroTik NMS deploy webhook agent
Documentation=https://github.com/c0de-ch/mikrotik-c0de
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=mikrotik-nms-agent
Group=mikrotik-nms-agent
WorkingDirectory=/var/lib/mikrotik-nms-deploy-agent
EnvironmentFile=/etc/mikrotik-nms-deploy-agent/env
ExecStart=/opt/mikrotik-nms-deploy-agent/bin/deploy-agent
Restart=on-failure
RestartSec=5
LimitNOFILE=4096

# Hardening
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
ProtectClock=true
ProtectHostname=true
ProtectProc=invisible
RestrictSUIDSGID=true
RestrictNamespaces=true
RestrictRealtime=true
LockPersonality=true
ReadWritePaths=/var/lib/mikrotik-nms-deploy-agent

[Install]
WantedBy=multi-user.target
